# debian/case-organizer.service
[Unit]
Description=Case Organizer (Flask)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/case-organizer
# use system Python, not a venv
ExecStart=/usr/bin/python3 /opt/case-organizer/app.py
# run under a dedicated user; created in postinst (see below)
User=caseorganizer
Group=caseorganizer
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

